# pkg_scripts/hello-world.conf
description "hello-world"

expect fork

start on runlevel [2345]
stop on runlevel [!2345]

exec java -jar /opt/hello-world/hello-world*.jar 2>&1 > /var/log/hello-world/hello.log &
